# docker/php/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; PHP-FPM - main application server
[program:php-fpm]
process_name=%(program_name)s
command=php-fpm -F -R
autostart=true
autorestart=true
stopasgroup=true
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/php-fpm.log
environment=LARAVEL_SKIP_ENV_FILE=true

; Queue worker - for processing queued jobs
; Adjust numprocs and processes based on your needs
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=false
autorestart=true
stopasgroup=true
stopwaitsecs=3600
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker.log
environment=LARAVEL_SKIP_ENV_FILE=true

; Scheduler - for running scheduled tasks
; The scheduler needs a cron trigger, but we run it continuously
[program:laravel-scheduler]
process_name=%(program_name)s
command=bash -c "while [ true ]; do php /var/www/html/artisan schedule:run >> /var/log/supervisor/scheduler.log 2>&1; sleep 60; done"
autostart=false
autorestart=true
stopasgroup=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/scheduler.log
environment=LARAVEL_SKIP_ENV_FILE=true

; Horizon - for monitoring queues (optional, if using Horizon)
[program:laravel-horizon]
process_name=%(program_name)s
command=php /var/www/html/artisan horizon
autostart=false
autorestart=true
stopasgroup=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/horizon.log
environment=LARAVEL_SKIP_ENV_FILE=true
